# install/ansible/roles/jukeberry/files/99-mount.rules
KERNEL=="sda", RUN+="/bin/mount /var/media/music/usb1"
KERNEL=="sdb", RUN+="/bin/mount /var/media/music/usb2"

